# drivemond-admin-new-install-3.1/deploy/supervisor/vito-reverb.conf
; Supervisor program for the Vito Reverb websocket server.
; Required in production: ride/mart chat and live tracking broadcast over Reverb
; (BROADCAST_DRIVER=reverb); without this process every private-channel
; subscription fails and the apps fall back to polling only.
;
; Topology (matches DEPLOY.md): Reverb listens on 127.0.0.1:8080; nginx
; terminates TLS on :6015 and proxies to it. Keep REVERB_PORT=8080 and
; PUSHER_PORT=6015 in .env in sync with this.
;
; Install:
;   sudo cp deploy/supervisor/vito-reverb.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start vito-reverb
;
; Adjust the paths/user below to your deployment.

[program:vito-reverb]
command=php /var/www/vito/artisan reverb:start --host=127.0.0.1 --port=8080 --no-interaction
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/vito-reverb.log
